# raspberry-pi/systemd/fallback-engine.service
[Unit]
Description=Smart Farm Fallback Engine (Emergency Failover)
Documentation=docs/02-design/features/rpi-emergency-failover.design.md
After=network-online.target chrony.service mosquitto.service gpio-agent.service
Wants=network-online.target

[Service]
Type=simple
# 최소 권한 원칙: pi 유저 + smartfarm 그룹.
# /var/lib/smartfarm/fallback 디렉터리는 setup.sh에서 pi:pi로 chown.
User=pi
Group=pi
# 무전압 우적센서(rain-gpio)가 /dev/gpiochip0(root:gpio) 를 gpioget 으로 읽으려면 gpio 그룹 필요.
# (프로덕션에서는 usermod -aG gpio pi 로도 부여됨 — fresh 이미지 재현성 위해 유닛에도 명시)
SupplementaryGroups=gpio
WorkingDirectory=/opt/smart-farm/fallback-engine
EnvironmentFile=/etc/smartfarm/fallback-engine.env
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=10
# SD 카드 쓰기 빈도 절감용 환경 변수 (event-queue 버퍼링)
Environment=FALLBACK_EVENT_BUFFER_MS=30000
Environment=FALLBACK_EVENT_BUFFER_MAX=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fallback-engine

[Install]
WantedBy=multi-user.target
