# raspberry-pi/systemd/zigbee2mqtt.service
[Unit]
Description=Zigbee2MQTT - Smart Farm Zigbee Gateway
After=network-online.target mosquitto.service
Wants=network-online.target
Requires=mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zigbee2mqtt
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# 환경 변수
Environment=NODE_ENV=production

# Watchdog (30초 내 응답 없으면 재시작)
WatchdogSec=30

[Install]
WantedBy=multi-user.target
